# packages/osquery_manager/kibana/osquery_saved_query/osquery_manager-f6a7b8c9-d0e1-34f0-5678-901234567890.json
{
    "attributes": {
        "created_at": "2026-01-27T00:00:00.000Z",
        "created_by": "elastic",
        "description": "Detects macOS persistence via dual-detection approach: (1) Non-Apple signed LaunchAgents/Daemons and (2) Living off the Land (LotL) attack indicators using legitimate macOS tools. Identifies both unsigned persistence AND abuse of bash, curl/osascript, base64, etc. Filters out Apple system paths while detecting suspicious patterns regardless of code signature. Note: macOS 10.13+ login items may not be fully captured due to Apple's backgrounditems.btm format. Maps to MITRE ATT&CK T1543.001 (Launch Agent/Daemon), T1547.015 (Login Items), T1059.004 (Unix Shell), T1105 (Ingress Tool Transfer).",
        "ecs_mapping": [
            {
                "key": "event.category",
                "value": {
                    "value": [
                        "configuration"
                    ]
                }
            },
            {
                "key": "event.type",
                "value": {
                    "value": [
                        "info"
                    ]
                }
            },
            {
                "key": "event.action",
                "value": {
                    "value": "osquery.startup_items"
                }
            },
            {
                "key": "process.executable",
                "value": {
                    "field": "exe_path"
                }
            },
            {
                "key": "process.command_line",
                "value": {
                    "field": "args"
                }
            },
            {
                "key": "file.path",
                "value": {
                    "field": "path"
                }
            },
            {
                "key": "file.hash.sha256",
                "value": {
                    "field": "sha256"
                }
            },
            {
                "key": "file.hash.sha1",
                "value": {
                    "field": "sha1"
                }
            },
            {
                "key": "file.hash.md5",
                "value": {
                    "field": "md5"
                }
            },
            {
                "key": "file.size",
                "value": {
                    "field": "size"
                }
            },
            {
                "key": "file.mtime",
                "value": {
                    "field": "modified_time"
                }
            },
            {
                "key": "file.ctime",
                "value": {
                    "field": "changed_time"
                }
            },
            {
                "key": "file.accessed",
                "value": {
                    "field": "accessed_time"
                }
            },
            {
                "key": "file.created",
                "value": {
                    "field": "created_time"
                }
            },
            {
                "key": "file.uid",
                "value": {
                    "field": "uid"
                }
            },
            {
                "key": "file.gid",
                "value": {
                    "field": "gid"
                }
            },
            {
                "key": "file.mode",
                "value": {
                    "field": "mode"
                }
            },
            {
                "key": "user.name",
                "value": {
                    "field": "username"
                }
            },
            {
                "key": "service.state",
                "value": {
                    "field": "status"
                }
            },
            {
                "key": "rule.category",
                "value": {
                    "field": "type"
                }
            },
            {
                "key": "file.code_signature.status",
                "value": {
                    "field": "signature_status"
                }
            },
            {
                "key": "file.code_signature.subject_name",
                "value": {
                    "field": "signature_signer"
                }
            },
            {
                "key": "rule.description",
                "value": {
                    "field": "detection_reason"
                }
            },
            {
                "key": "rule.name",
                "value": {
                    "field": "detection_method"
                }
            },
            {
                "key": "tags",
                "value": {
                    "value": [
                        "osquery",
                        "persistence",
                        "startup_items",
                        "macos"
                    ]
                }
            }
        ],
        "id": "startup_items_darwin_elastic",
        "interval": "3600",
        "platform": "darwin",
        "timeout": 180,
        "query": "-- Dual-detection macOS persistence query:\n-- 1. NON_WHITELISTED: Non-Apple signed LaunchAgents/Daemons (signature-based filtering)\n-- 2. LOTL_INDICATOR: Living off the Land patterns (bash -c, curl | bash, osascript, etc.)\n-- Uses TRIM() on paths and extracts executable path from program_arguments\n-- MITRE ATT&CK: T1543.001, T1547.015, T1059.004, T1105\n\nWITH non_whitelisted_launchd AS (\n    SELECT\n        l.label AS name,\n        TRIM(COALESCE(NULLIF(l.program, ''), l.program_arguments)) AS path,\n        l.path AS source,\n        l.program_arguments AS args,\n        l.username,\n        CASE\n            WHEN l.disabled = '1' OR l.disabled = 1 THEN 'disabled'\n            WHEN l.run_at_load = 'true' OR l.run_at_load = '1' OR l.run_at_load = 1 THEN 'enabled'\n            ELSE 'unknown'\n        END AS status,\n        'Launch Agent/Daemon' AS type,\n        'NON_WHITELISTED' AS detection_method,\n        'Non-Apple signed LaunchAgent/Daemon' AS detection_reason\n    FROM launchd AS l\n    WHERE (\n        (l.program IS NOT NULL AND TRIM(l.program) != '')\n        OR (l.program_arguments IS NOT NULL AND TRIM(l.program_arguments) != '')\n    )\n        AND l.path NOT LIKE '/System/Library/%'\n        AND l.path NOT LIKE '/Library/Apple/%'\n        AND (l.run_at_load = 'true' OR l.run_at_load = '1' OR l.run_at_load = 1)\n),\nnon_whitelisted_startup AS (\n    SELECT\n        si.name,\n        TRIM(si.path) AS path,\n        si.source,\n        si.args,\n        si.username,\n        si.status,\n        CASE\n            WHEN si.type = 'Startup Item' THEN 'Legacy Startup Item'\n            WHEN si.type = 'Login Item' THEN 'Login Item'\n            ELSE si.type\n        END AS type,\n        'NON_WHITELISTED' AS detection_method,\n        'Legacy startup/login item' AS detection_reason\n    FROM startup_items AS si\n    WHERE TRIM(si.path) IS NOT NULL\n        AND TRIM(si.path) != ''\n        AND TRIM(si.path) NOT LIKE '/System/Library/%'\n        AND TRIM(si.path) NOT LIKE '/Library/Apple/%'\n),\nlotl_launchd AS (\n    SELECT\n        l.label AS name,\n        TRIM(COALESCE(NULLIF(l.program, ''), l.program_arguments)) AS path,\n        l.path AS source,\n        l.program_arguments AS args,\n        l.username,\n        CASE\n            WHEN l.disabled = '1' OR l.disabled = 1 THEN 'disabled'\n            WHEN l.run_at_load = 'true' OR l.run_at_load = '1' OR l.run_at_load = 1 THEN 'enabled'\n            ELSE 'unknown'\n        END AS status,\n        'Launch Agent/Daemon (LotL)' AS type,\n        'LOTL_INDICATOR' AS detection_method,\n        CASE\n            WHEN l.program_arguments LIKE '%bash -c%' OR l.program_arguments LIKE '%sh -c%' THEN 'Shell command execution via -c flag'\n            WHEN l.program_arguments LIKE '%curl%|%bash%' OR l.program_arguments LIKE '%curl%|%sh%' THEN 'Download and pipe to shell'\n            WHEN l.program_arguments LIKE '%curl%http%' OR l.program_arguments LIKE '%wget%http%' THEN 'Download utility abuse'\n            WHEN l.program_arguments LIKE '%base64 -D%' OR l.program_arguments LIKE '%base64 --decode%' THEN 'Base64 decode for obfuscation'\n            WHEN l.program_arguments LIKE '%osascript%' AND l.program_arguments LIKE '%-e%' THEN 'AppleScript execution abuse'\n            WHEN l.program_arguments LIKE '%python%-%c%' OR l.program_arguments LIKE '%perl%-%e%' THEN 'Scripting language one-liner'\n            WHEN (l.program_arguments LIKE '%/tmp/%' AND l.program_arguments NOT LIKE '%/var/tmp/%') OR l.program_arguments LIKE '%/private/tmp/%' THEN 'Execution from temp directory'\n            WHEN l.program_arguments LIKE '%/dev/tcp/%' THEN 'Bash TCP socket redirection'\n            WHEN l.program_arguments LIKE '% nc %' OR l.program_arguments LIKE '%/nc %' OR l.program_arguments LIKE 'nc %' OR l.program_arguments LIKE '%netcat%' OR l.program_arguments LIKE '%ncat %' THEN 'Netcat reverse shell'\n            WHEN l.program_arguments LIKE '%nohup%&%' OR l.program_arguments LIKE '%disown%' THEN 'Background process persistence'\n            WHEN l.program_arguments LIKE '%.sh%' AND (TRIM(COALESCE(NULLIF(l.program, ''), l.program_arguments)) LIKE '/tmp/%' OR TRIM(COALESCE(NULLIF(l.program, ''), l.program_arguments)) LIKE '/private/tmp/%') THEN 'Shell script from temp directory'\n            ELSE 'Unknown LotL pattern in LaunchAgent'\n        END AS detection_reason\n    FROM launchd AS l\n    WHERE (\n        (l.program IS NOT NULL AND TRIM(l.program) != '')\n        OR (l.program_arguments IS NOT NULL AND TRIM(l.program_arguments) != '')\n    )\n        AND (l.run_at_load = 'true' OR l.run_at_load = '1' OR l.run_at_load = 1)\n        AND (\n            l.program_arguments LIKE '%bash -c%'\n            OR l.program_arguments LIKE '%sh -c%'\n            OR l.program_arguments LIKE '%curl%|%bash%'\n            OR l.program_arguments LIKE '%curl%|%sh%'\n            OR l.program_arguments LIKE '%curl%http%'\n            OR l.program_arguments LIKE '%wget%http%'\n            OR l.program_arguments LIKE '%base64 -D%'\n            OR l.program_arguments LIKE '%base64 --decode%'\n            OR (l.program_arguments LIKE '%osascript%' AND l.program_arguments LIKE '%-e%')\n            OR l.program_arguments LIKE '%python%-%c%'\n            OR l.program_arguments LIKE '%perl%-%e%'\n            OR (l.program_arguments LIKE '%/tmp/%' AND l.program_arguments NOT LIKE '%/var/tmp/%')\n            OR l.program_arguments LIKE '%/private/tmp/%'\n            OR l.program_arguments LIKE '%/dev/tcp/%'\n            OR (l.program_arguments LIKE '% nc %' OR l.program_arguments LIKE '%/nc %' OR l.program_arguments LIKE 'nc %')\n            OR l.program_arguments LIKE '%netcat%'\n            OR l.program_arguments LIKE '%ncat %'\n            OR l.program_arguments LIKE '%nohup%&%'\n            OR l.program_arguments LIKE '%disown%'\n            OR (l.program_arguments LIKE '%.sh%' AND (TRIM(COALESCE(NULLIF(l.program, ''), l.program_arguments)) LIKE '/tmp/%' OR TRIM(COALESCE(NULLIF(l.program, ''), l.program_arguments)) LIKE '/private/tmp/%'))\n        )\n),\ncombined AS (\n    SELECT * FROM non_whitelisted_launchd\n    UNION ALL\n    SELECT * FROM non_whitelisted_startup\n    UNION ALL\n    SELECT * FROM lotl_launchd\n)\nSELECT\n    c.name,\n    c.path,\n    -- Extract executable path (first space-separated token if path contains arguments)\n    CASE\n        WHEN INSTR(c.path, ' ') > 0 AND SUBSTR(c.path, 1, 1) != '-'\n        THEN SUBSTR(c.path, 1, INSTR(c.path, ' ') - 1)\n        ELSE c.path\n    END AS exe_path,\n    c.type,\n    c.status,\n    c.source,\n    c.args,\n    c.username,\n    c.detection_method,\n    c.detection_reason,\n    CASE WHEN s.signed = 1 THEN 'signed' ELSE 'unsigned' END AS signature_status,\n    s.identifier AS signature_signer,\n    h.sha256,\n    h.sha1,\n    h.md5,\n    f.size,\n    datetime(f.mtime, 'unixepoch') AS modified_time,\n    datetime(f.ctime, 'unixepoch') AS changed_time,\n    datetime(f.atime, 'unixepoch') AS accessed_time,\n    datetime(f.btime, 'unixepoch') AS created_time,\n    f.uid,\n    f.gid,\n    f.mode\nFROM combined AS c\nLEFT JOIN signature AS s ON s.path = CASE\n    WHEN INSTR(c.path, ' ') > 0 AND SUBSTR(c.path, 1, 1) != '-'\n    THEN SUBSTR(c.path, 1, INSTR(c.path, ' ') - 1)\n    ELSE c.path\nEND\nLEFT JOIN hash AS h ON h.path = CASE\n    WHEN INSTR(c.path, ' ') > 0 AND SUBSTR(c.path, 1, 1) != '-'\n    THEN SUBSTR(c.path, 1, INSTR(c.path, ' ') - 1)\n    ELSE c.path\nEND\nLEFT JOIN file AS f ON f.path = CASE\n    WHEN INSTR(c.path, ' ') > 0 AND SUBSTR(c.path, 1, 1) != '-'\n    THEN SUBSTR(c.path, 1, INSTR(c.path, ' ') - 1)\n    ELSE c.path\nEND\nWHERE (\n    c.detection_method = 'LOTL_INDICATOR'\n    OR s.signed IS NULL\n    OR s.signed = 0\n    OR (\n        s.identifier IS NOT NULL\n        AND s.identifier NOT LIKE 'com.apple.%'\n        AND s.identifier NOT LIKE 'Apple Inc.%'\n    )\n)\nORDER BY\n    CASE WHEN c.detection_method = 'LOTL_INDICATOR' THEN 0 ELSE 1 END,\n    CASE WHEN s.signed = 0 OR s.signed IS NULL THEN 0 ELSE 1 END,\n    c.detection_reason,\n    c.name",
        "updated_at": "2026-01-27T00:00:00.000Z",
        "updated_by": "elastic"
    },
    "coreMigrationVersion": "9.3.0",
    "id": "osquery_manager-f6a7b8c9-d0e1-34f0-5678-901234567890",
    "references": [],
    "type": "osquery-saved-query",
    "updated_at": "2026-01-27T00:00:00.000Z",
    "version": "WzEsMV0="
}
