# packages/osquery_manager/kibana/osquery_saved_query/osquery_manager-892ee425-60e7-4eb6-ba25-6e97dc3e2ea0.json
{
  "attributes": {
    "created_at": "2025-11-05T09:19:01.000Z",
    "created_by": "elastic",
    "description": "Detect suspicious Windows services with comprehensive risk scoring (0-100). Identifies: (1) Services in user-writable directories (Temp, AppData, Downloads), (2) Unsigned or untrusted binaries via code signature validation, (3) ServiceDLL hijacking (missing DLLs, unsigned DLLs, suspicious locations), (4) Failure command persistence mechanisms, (5) Privilege escalation (SYSTEM account with user-writable paths). Risk Scoring: CRITICAL (70-100) - Immediate investigation required, HIGH (50-69) - Review within 24 hours, MEDIUM (30-49) - Weekly review, LOW (0-29) - Baseline monitoring. Covers MITRE ATT&CK: T1543.003 (Create or Modify System Process: Windows Service), T1574.011 (Hijack Execution Flow: Services Registry Permissions Weakness)",
    "ecs_mapping": [
      {
        "key": "service.name",
        "value": {
          "field": "name"
        }
      },
      {
        "key": "service.state",
        "value": {
          "field": "status"
        }
      },
      {
        "key": "process.executable",
        "value": {
          "field": "path"
        }
      },
      {
        "key": "user.name",
        "value": {
          "field": "user_account"
        }
      },
      {
        "key": "file.hash.md5",
        "value": {
          "field": "md5"
        }
      },
      {
        "key": "file.hash.sha256",
        "value": {
          "field": "sha256"
        }
      },
      {
        "key": "file.code_signature.subject_name",
        "value": {
          "field": "cert_subject"
        }
      },
      {
        "key": "file.code_signature.status",
        "value": {
          "field": "signature_status"
        }
      },
      {
        "key": "dll.hash.sha256",
        "value": {
          "field": "servicedll_sha256"
        }
      },
      {
        "key": "dll.hash.md5",
        "value": {
          "field": "servicedll_md5"
        }
      },
      {
        "key": "dll.code_signature.subject_name",
        "value": {
          "field": "servicedll_cert_subject"
        }
      },
      {
        "key": "dll.code_signature.status",
        "value": {
          "field": "servicedll_signature_status"
        }
      },
      {
        "key": "event.risk_score",
        "value": {
          "field": "risk_score"
        }
      },
      {
        "key": "event.severity",
        "value": {
          "field": "risk_level"
        }
      },
      {
        "key": "event.created",
        "value": {
          "field": "service_registry_mtime"
        }
      },
      {
        "key": "event.category",
        "value": {
          "value": [
            "configuration"
          ]
        }
      },
      {
        "key": "event.type",
        "value": {
          "value": [
            "info"
          ]
        }
      },
      {
        "key": "tags",
        "value": {
          "value": [
            "persistence",
            "services",
            "backdoor_detection",
            "service_dll",
            "risk_scoring",
            "servicedll_hijacking",
            "privilege_escalation",
            "code_signing"
          ]
        }
      }
    ],
    "id": "suspicious_services_elastic",
    "interval": "3600",
    "platform": "windows",
    "query": "SELECT \n  s.name, \n  s.display_name, \n  s.status, \n  s.start_type, \n  s.path, \n  s.user_account,\n  r1.data AS service_dll,\n  r2.data AS failure_command,\n  h.sha256,\n  h.md5,\n  a.subject_name AS cert_subject,\n  a.issuer_name AS cert_issuer,\n  a.result AS signature_status,\n  h_dll.sha256 AS servicedll_sha256,\n  h_dll.md5 AS servicedll_md5,\n  a_dll.subject_name AS servicedll_cert_subject,\n  a_dll.result AS servicedll_signature_status,\n  r_service_key.mtime AS service_registry_mtime,\n  CAST((strftime('%s', 'now') - r_service_key.mtime) / 86400 AS INTEGER) AS days_since_creation,\n  (\n    (CASE\n      WHEN s.path LIKE '%\\\\Temp\\\\%' THEN 40\n      WHEN s.path LIKE '%\\\\Users\\\\%\\\\AppData\\\\Local\\\\Temp\\\\%' THEN 40\n      WHEN s.path LIKE '%\\\\Users\\\\%\\\\Downloads\\\\%' THEN 35\n      WHEN s.path LIKE '%\\\\AppData\\\\Local\\\\%' THEN 35\n      WHEN s.path LIKE '%\\\\AppData\\\\Roaming\\\\%' THEN 30\n      WHEN s.path LIKE '%\\\\Users\\\\%\\\\Documents\\\\%' THEN 30\n      WHEN s.path LIKE '%\\\\Users\\\\Public\\\\%' THEN 25\n      WHEN s.path LIKE '%\\\\ProgramData\\\\%' AND s.path NOT LIKE '%\\\\ProgramData\\\\Microsoft\\\\%' THEN 20\n      ELSE 5\n    END) +\n    (CASE\n      WHEN a.result IS NULL OR a.result = '' THEN 25\n      WHEN a.result != 'trusted' THEN 25\n      WHEN a.subject_name NOT LIKE '%Microsoft%' AND a.subject_name NOT LIKE '%Windows%' THEN 15\n      ELSE 0\n    END) +\n    (CASE\n      WHEN s.user_account IN ('LocalSystem', 'NT AUTHORITY\\\\SYSTEM', 'SYSTEM') AND\n           (s.path LIKE '%\\\\Temp\\\\%' OR s.path LIKE '%\\\\AppData\\\\%' OR s.path LIKE '%\\\\Users\\\\%') THEN 20\n      WHEN s.user_account IN ('LocalSystem', 'NT AUTHORITY\\\\SYSTEM', 'SYSTEM') THEN 10\n      WHEN s.user_account LIKE '%Administrator%' THEN 5\n      ELSE 0\n    END) +\n    (CASE\n      WHEN LOWER(s.name) IN ('service', 'svc', 'system', 'update', 'windows', 'microsoft', 'security', 'monitor', 'agent', 'helper') THEN 10\n      WHEN LENGTH(s.name) <= 5 THEN 8\n      WHEN s.display_name = s.name OR s.display_name IS NULL THEN 5\n      ELSE 0\n    END) +\n    (CASE\n      WHEN r1.data IS NOT NULL AND r1.data != '' AND h_dll.sha256 IS NULL THEN 5\n      WHEN r1.data LIKE '%\\\\Temp\\\\%' OR r1.data LIKE '%\\\\AppData\\\\%' THEN 5\n      WHEN a_dll.result IS NOT NULL AND a_dll.result != 'trusted' THEN 5\n      ELSE 0\n    END)\n  ) AS risk_score,\n  (CASE\n    WHEN (\n      (CASE\n        WHEN s.path LIKE '%\\\\Temp\\\\%' THEN 40\n        WHEN s.path LIKE '%\\\\Users\\\\%\\\\AppData\\\\Local\\\\Temp\\\\%' THEN 40\n        WHEN s.path LIKE '%\\\\Users\\\\%\\\\Downloads\\\\%' THEN 35\n        WHEN s.path LIKE '%\\\\AppData\\\\Local\\\\%' THEN 35\n        WHEN s.path LIKE '%\\\\AppData\\\\Roaming\\\\%' THEN 30\n        WHEN s.path LIKE '%\\\\Users\\\\%\\\\Documents\\\\%' THEN 30\n        WHEN s.path LIKE '%\\\\Users\\\\Public\\\\%' THEN 25\n        WHEN s.path LIKE '%\\\\ProgramData\\\\%' AND s.path NOT LIKE '%\\\\ProgramData\\\\Microsoft\\\\%' THEN 20\n        ELSE 5\n      END) +\n      (CASE\n        WHEN a.result IS NULL OR a.result = '' THEN 25\n        WHEN a.result != 'trusted' THEN 25\n        WHEN a.subject_name NOT LIKE '%Microsoft%' AND a.subject_name NOT LIKE '%Windows%' THEN 15\n        ELSE 0\n      END) +\n      (CASE\n        WHEN s.user_account IN ('LocalSystem', 'NT AUTHORITY\\\\SYSTEM', 'SYSTEM') AND\n             (s.path LIKE '%\\\\Temp\\\\%' OR s.path LIKE '%\\\\AppData\\\\%' OR s.path LIKE '%\\\\Users\\\\%') THEN 20\n        WHEN s.user_account IN ('LocalSystem', 'NT AUTHORITY\\\\SYSTEM', 'SYSTEM') THEN 10\n        WHEN s.user_account LIKE '%Administrator%' THEN 5\n        ELSE 0\n      END) +\n      (CASE\n        WHEN LOWER(s.name) IN ('service', 'svc', 'system', 'update', 'windows', 'microsoft', 'security', 'monitor', 'agent', 'helper') THEN 10\n        WHEN LENGTH(s.name) <= 5 THEN 8\n        WHEN s.display_name = s.name OR s.display_name IS NULL THEN 5\n        ELSE 0\n      END) +\n      (CASE\n        WHEN r1.data IS NOT NULL AND r1.data != '' AND h_dll.sha256 IS NULL THEN 5\n        WHEN r1.data LIKE '%\\\\Temp\\\\%' OR r1.data LIKE '%\\\\AppData\\\\%' THEN 5\n        WHEN a_dll.result IS NOT NULL AND a_dll.result != 'trusted' THEN 5\n        ELSE 0\n      END)\n    ) >= 70 THEN 'CRITICAL'\n    WHEN (\n      (CASE\n        WHEN s.path LIKE '%\\\\Temp\\\\%' THEN 40\n        WHEN s.path LIKE '%\\\\Users\\\\%\\\\AppData\\\\Local\\\\Temp\\\\%' THEN 40\n        WHEN s.path LIKE '%\\\\Users\\\\%\\\\Downloads\\\\%' THEN 35\n        WHEN s.path LIKE '%\\\\AppData\\\\Local\\\\%' THEN 35\n        WHEN s.path LIKE '%\\\\AppData\\\\Roaming\\\\%' THEN 30\n        WHEN s.path LIKE '%\\\\Users\\\\%\\\\Documents\\\\%' THEN 30\n        WHEN s.path LIKE '%\\\\Users\\\\Public\\\\%' THEN 25\n        WHEN s.path LIKE '%\\\\ProgramData\\\\%' AND s.path NOT LIKE '%\\\\ProgramData\\\\Microsoft\\\\%' THEN 20\n        ELSE 5\n      END) +\n      (CASE\n        WHEN a.result IS NULL OR a.result = '' THEN 25\n        WHEN a.result != 'trusted' THEN 25\n        WHEN a.subject_name NOT LIKE '%Microsoft%' AND a.subject_name NOT LIKE '%Windows%' THEN 15\n        ELSE 0\n      END) +\n      (CASE\n        WHEN s.user_account IN ('LocalSystem', 'NT AUTHORITY\\\\SYSTEM', 'SYSTEM') AND\n             (s.path LIKE '%\\\\Temp\\\\%' OR s.path LIKE '%\\\\AppData\\\\%' OR s.path LIKE '%\\\\Users\\\\%') THEN 20\n        WHEN s.user_account IN ('LocalSystem', 'NT AUTHORITY\\\\SYSTEM', 'SYSTEM') THEN 10\n        WHEN s.user_account LIKE '%Administrator%' THEN 5\n        ELSE 0\n      END) +\n      (CASE\n        WHEN LOWER(s.name) IN ('service', 'svc', 'system', 'update', 'windows', 'microsoft', 'security', 'monitor', 'agent', 'helper') THEN 10\n        WHEN LENGTH(s.name) <= 5 THEN 8\n        WHEN s.display_name = s.name OR s.display_name IS NULL THEN 5\n        ELSE 0\n      END) +\n      (CASE\n        WHEN r1.data IS NOT NULL AND r1.data != '' AND h_dll.sha256 IS NULL THEN 5\n        WHEN r1.data LIKE '%\\\\Temp\\\\%' OR r1.data LIKE '%\\\\AppData\\\\%' THEN 5\n        WHEN a_dll.result IS NOT NULL AND a_dll.result != 'trusted' THEN 5\n        ELSE 0\n      END)\n    ) >= 50 THEN 'HIGH'\n    WHEN (\n      (CASE\n        WHEN s.path LIKE '%\\\\Temp\\\\%' THEN 40\n        WHEN s.path LIKE '%\\\\Users\\\\%\\\\AppData\\\\Local\\\\Temp\\\\%' THEN 40\n        WHEN s.path LIKE '%\\\\Users\\\\%\\\\Downloads\\\\%' THEN 35\n        WHEN s.path LIKE '%\\\\AppData\\\\Local\\\\%' THEN 35\n        WHEN s.path LIKE '%\\\\AppData\\\\Roaming\\\\%' THEN 30\n        WHEN s.path LIKE '%\\\\Users\\\\%\\\\Documents\\\\%' THEN 30\n        WHEN s.path LIKE '%\\\\Users\\\\Public\\\\%' THEN 25\n        WHEN s.path LIKE '%\\\\ProgramData\\\\%' AND s.path NOT LIKE '%\\\\ProgramData\\\\Microsoft\\\\%' THEN 20\n        ELSE 5\n      END) +\n      (CASE\n        WHEN a.result IS NULL OR a.result = '' THEN 25\n        WHEN a.result != 'trusted' THEN 25\n        WHEN a.subject_name NOT LIKE '%Microsoft%' AND a.subject_name NOT LIKE '%Windows%' THEN 15\n        ELSE 0\n      END) +\n      (CASE\n        WHEN s.user_account IN ('LocalSystem', 'NT AUTHORITY\\\\SYSTEM', 'SYSTEM') AND\n             (s.path LIKE '%\\\\Temp\\\\%' OR s.path LIKE '%\\\\AppData\\\\%' OR s.path LIKE '%\\\\Users\\\\%') THEN 20\n        WHEN s.user_account IN ('LocalSystem', 'NT AUTHORITY\\\\SYSTEM', 'SYSTEM') THEN 10\n        WHEN s.user_account LIKE '%Administrator%' THEN 5\n        ELSE 0\n      END) +\n      (CASE\n        WHEN LOWER(s.name) IN ('service', 'svc', 'system', 'update', 'windows', 'microsoft', 'security', 'monitor', 'agent', 'helper') THEN 10\n        WHEN LENGTH(s.name) <= 5 THEN 8\n        WHEN s.display_name = s.name OR s.display_name IS NULL THEN 5\n        ELSE 0\n      END) +\n      (CASE\n        WHEN r1.data IS NOT NULL AND r1.data != '' AND h_dll.sha256 IS NULL THEN 5\n        WHEN r1.data LIKE '%\\\\Temp\\\\%' OR r1.data LIKE '%\\\\AppData\\\\%' THEN 5\n        WHEN a_dll.result IS NOT NULL AND a_dll.result != 'trusted' THEN 5\n        ELSE 0\n      END)\n    ) >= 30 THEN 'MEDIUM'\n    ELSE 'LOW'\n  END) AS risk_level\nFROM services s\nLEFT JOIN registry r1 ON r1.path = 'HKEY_LOCAL_MACHINE\\\\System\\\\CurrentControlSet\\\\Services\\\\' || s.name || '\\\\Parameters' AND r1.name = 'ServiceDll'\nLEFT JOIN registry r2 ON r2.path = 'HKEY_LOCAL_MACHINE\\\\System\\\\CurrentControlSet\\\\Services\\\\' || s.name AND r2.name = 'FailureCommand'\nLEFT JOIN hash h ON h.path = s.path\nLEFT JOIN authenticode a ON a.path = s.path\nLEFT JOIN hash h_dll ON h_dll.path = r1.data\nLEFT JOIN authenticode a_dll ON a_dll.path = r1.data\nLEFT JOIN registry r_service_key ON r_service_key.path = 'HKEY_LOCAL_MACHINE\\\\System\\\\CurrentControlSet\\\\Services\\\\' || s.name\nWHERE\n  s.start_type IN ('AUTO_START', 'DEMAND_START', 'BOOT_START', 'SYSTEM_START')\n  AND (\n    s.path LIKE '%\\\\Temp\\\\%'\n    OR s.path LIKE '%\\\\Users\\\\%\\\\AppData\\\\%'\n    OR s.path LIKE '%\\\\Users\\\\%\\\\Documents\\\\%'\n    OR s.path LIKE '%\\\\Users\\\\%\\\\Downloads\\\\%'\n    OR s.path LIKE '%\\\\Users\\\\Public\\\\%'\n    OR (s.path LIKE '%\\\\ProgramData\\\\%' AND s.path NOT LIKE '%\\\\ProgramData\\\\Microsoft\\\\%')\n    OR (r1.data IS NOT NULL AND r1.data != '' AND r1.data NOT LIKE 'C:\\\\Windows\\\\System32\\\\%' AND r1.data NOT LIKE 'C:\\\\Windows\\\\SysWOW64\\\\%')\n    OR (r2.data IS NOT NULL AND r2.data != '')\n    OR (a.result IS NOT NULL AND a.result != 'trusted')\n  )\n  AND s.path NOT LIKE 'C:\\\\Windows\\\\System32\\\\%'\n  AND s.path NOT LIKE 'C:\\\\Windows\\\\SysWOW64\\\\%'\n  AND s.name NOT IN ('wuauserv', 'BITS', 'TrustedInstaller', 'WinDefend', 'Spooler')\nORDER BY risk_score DESC\nLIMIT 200;",
    "updated_at": "2025-11-05T09:19:01.000Z",
    "updated_by": "elastic"
  },
  "coreMigrationVersion": "9.2.0",
  "id": "osquery_manager-892ee425-60e7-4eb6-ba25-6e97dc3e2ea0",
  "references": [],
  "type": "osquery-saved-query",
  "updated_at": "2025-11-05T09:19:01.000Z",
  "version": "WzgzLDJd"
}
