# packages/osquery_manager/kibana/osquery_saved_query/osquery_manager-f8b0894b-772d-4242-8e19-dbc5d7ae2e06.json
{
    "attributes": {
        "created_at": "2026-01-27T00:00:00.000Z",
        "created_by": "elastic",
        "description": "Identifies suspicious Linux systemd service units in non-standard locations. Targets user-writable paths (/home, /tmp, /var/tmp), user systemd directories (~/.config/systemd, ~/.local/share/systemd), unusual unit states (linked, transient), recently modified units in /etc/systemd/system, and runtime-generated units in /run. Prioritizes results by risk level.",
        "ecs_mapping": [
            {
                "key": "event.category",
                "value": {
                    "value": [
                        "configuration"
                    ]
                }
            },
            {
                "key": "event.type",
                "value": {
                    "value": [
                        "info"
                    ]
                }
            },
            {
                "key": "event.action",
                "value": {
                    "value": "osquery.systemd_services"
                }
            },
            {
                "key": "service.name",
                "value": {
                    "field": "unit_id"
                }
            },
            {
                "key": "service.state",
                "value": {
                    "field": "active_state"
                }
            },
            {
                "key": "file.path",
                "value": {
                    "field": "fragment_path"
                }
            },
            {
                "key": "file.mtime",
                "value": {
                    "field": "modified_time"
                }
            },
            {
                "key": "tags",
                "value": {
                    "value": [
                        "osquery",
                        "persistence",
                        "threat_hunting",
                        "systemd",
                        "linux"
                    ]
                }
            }
        ],
        "id": "services_suspicious_linux_elastic",
        "interval": "3600",
        "timeout": 120,
        "platform": "linux",
        "query": "-- Linux Systemd Service Persistence Detection\n-- Identifies suspicious service units in non-standard locations\n\nWITH services_base AS (\n  SELECT\n    su.id AS unit_id,\n    su.load_state,\n    su.active_state,\n    su.sub_state,\n    su.description,\n    su.unit_file_state,\n    su.fragment_path,\n    datetime(f.mtime, 'unixepoch') AS modified_time,\n    CAST((strftime('%s', 'now') - COALESCE(f.mtime, 0)) / 86400 AS INTEGER) AS days_since_modified\n  FROM systemd_units su\n  LEFT JOIN file f ON f.path = su.fragment_path\n  WHERE su.id LIKE '%.service'\n    AND su.fragment_path IS NOT NULL\n    AND su.fragment_path != ''\n)\n\nSELECT *\nFROM services_base\nWHERE (\n  fragment_path LIKE '/home/%'\n  OR fragment_path LIKE '/root/%'\n  OR fragment_path LIKE '/tmp/%'\n  OR fragment_path LIKE '/var/tmp/%'\n  OR fragment_path LIKE '%/.config/systemd/%'\n  OR fragment_path LIKE '%/.local/share/systemd/%'\n  OR unit_file_state IN ('linked', 'linked-runtime', 'transient')\n  OR (fragment_path LIKE '/etc/systemd/system/%' AND days_since_modified < 30)\n  OR fragment_path LIKE '/run/systemd/%'\n)\nORDER BY\n  CASE WHEN fragment_path LIKE '/home/%' OR fragment_path LIKE '/tmp/%' THEN 0 ELSE 1 END,\n  CASE WHEN fragment_path LIKE '%/.config/systemd/%' THEN 0 ELSE 1 END,\n  days_since_modified ASC;",
        "updated_at": "2026-01-27T00:00:00.000Z",
        "updated_by": "elastic"
    },
    "coreMigrationVersion": "9.3.0",
    "id": "osquery_manager-f8b0894b-772d-4242-8e19-dbc5d7ae2e06",
    "references": [],
    "type": "osquery-saved-query",
    "updated_at": "2026-01-27T00:00:00.000Z",
    "version": "WzEsMV0="
}
