# packages/osquery_manager/kibana/osquery_saved_query/osquery_manager-a1b2c3d4-5678-4abc-def0-123456789035.json
{
    "attributes": {
        "created_at": "2025-01-30T12:00:00.000Z",
        "created_by": "elastic",
        "description": "Detect macOS user-level persistence via LaunchAgents by monitoring for auto-start configurations, keep-alive settings, and suspicious program locations",
        "ecs_mapping": [
            {
                "key": "process.name",
                "value": {
                    "field": "name"
                }
            },
            {
                "key": "process.executable",
                "value": {
                    "field": "program"
                }
            },
            {
                "key": "process.args",
                "value": {
                    "field": "program_arguments"
                }
            },
            {
                "key": "file.path",
                "value": {
                    "field": "path"
                }
            },
            {
                "key": "event.category",
                "value": {
                    "value": ["process"]
                }
            },
            {
                "key": "event.type",
                "value": {
                    "value": ["info"]
                }
            },
            {
                "key": "event.action",
                "value": {
                    "value": "launchagent-enumeration"
                }
            },
            {
                "key": "event.dataset",
                "value": {
                    "value": "macos.launchd"
                }
            },
            {
                "key": "tags",
                "value": {
                    "value": ["persistence", "macos", "launchagent", "startup", "T1543.001"]
                }
            }
        ],
        "id": "macos_launchagents_persistence_elastic",
        "interval": "3600",
        "platform": "darwin",
        "query": "SELECT name, label, path, program, program_arguments, run_at_load, keep_alive, inetd_compatibility, start_on_mount, disabled FROM launchd WHERE path LIKE '%/LaunchAgents/%' AND (run_at_load = 1 OR keep_alive = 1 OR start_on_mount = 1 OR program LIKE '/tmp/%' OR program LIKE '%/Library/Application Support/%' OR program LIKE '%curl%' OR program LIKE '%bash -c%' OR program LIKE '%python%' OR program_arguments LIKE '%curl%' OR program_arguments LIKE '%wget%') ORDER BY disabled ASC, run_at_load DESC;",
        "updated_at": "2025-01-30T12:00:00.000Z",
        "updated_by": "elastic"
    },
    "coreMigrationVersion": "8.3.0",
    "id": "osquery_manager-a1b2c3d4-5678-4abc-def0-123456789035",
    "references": [],
    "type": "osquery-saved-query",
    "updated_at": "2025-01-30T12:00:00.000Z",
    "version": "WzEsMV0="
}
