# packages/osquery_manager/kibana/osquery_saved_query/osquery_manager-239dce60-b4a9-11ec-8f39-bf9c07530bbb.json
{
    "attributes": {
        "created_at": "2022-04-05T06:25:08.166Z",
        "created_by": "elastic",
        "description": "Get applications and executables configured to launch when a system reboots on Windows.",
        "ecs_mapping": [
            {
                "key": "file.path",
                "value": {
                    "field": "path"
                }
            },
            {
                "key": "file.name",
                "value": {
                    "field": "name"
                }
            },
            {
                "key": "process.args",
                "value": {
                    "field": "args"
                }
            },
            {
                "key": "file.directory",
                "value": {
                    "field": "source"
                }
            },
            {
                "key": "user.name",
                "value": {
                    "field": "username"
                }
            }
        ],
        "id": "persisted_apps_executables_windows_elastic",
        "interval": "3600",
        "platform": "windows",
        "query": "SELECT * FROM autoexec, startup_items;",
        "updated_at": "2022-04-05T06:25:08.166Z",
        "updated_by": "elastic"
    },
    "coreMigrationVersion": "8.3.0",
    "id": "osquery_manager-239dce60-b4a9-11ec-8f39-bf9c07530bbb",
    "references": [],
    "type": "osquery-saved-query",
    "updated_at": "2022-04-05T06:25:08.167Z",
    "version": "Wzc1MCwxXQ=="
}