# packages/osquery_manager/kibana/osquery_saved_query/osquery_manager-a1b2c3d4-5678-4abc-def0-123456789036.json
{
    "attributes": {
        "created_at": "2025-01-30T12:00:00.000Z",
        "created_by": "elastic",
        "description": "Detect macOS system-level persistence and privilege escalation via LaunchDaemons by monitoring for suspicious program locations, user directory access, and network listeners",
        "ecs_mapping": [
            {
                "key": "process.name",
                "value": {
                    "field": "name"
                }
            },
            {
                "key": "process.executable",
                "value": {
                    "field": "program"
                }
            },
            {
                "key": "process.args",
                "value": {
                    "field": "program_arguments"
                }
            },
            {
                "key": "file.path",
                "value": {
                    "field": "path"
                }
            },
            {
                "key": "event.category",
                "value": {
                    "value": ["process"]
                }
            },
            {
                "key": "event.type",
                "value": {
                    "value": ["info"]
                }
            },
            {
                "key": "event.action",
                "value": {
                    "value": "launchdaemon-enumeration"
                }
            },
            {
                "key": "event.dataset",
                "value": {
                    "value": "macos.launchd"
                }
            },
            {
                "key": "tags",
                "value": {
                    "value": ["persistence", "privilege_escalation", "macos", "launchdaemon", "system", "T1543.004"]
                }
            }
        ],
        "id": "macos_launchdaemons_persistence_elastic",
        "interval": "3600",
        "platform": "darwin",
        "query": "SELECT name, label, path, program, program_arguments, run_at_load, keep_alive, inetd_compatibility, start_on_mount, disabled FROM launchd WHERE path LIKE '%/LaunchDaemons/%' AND (run_at_load = 1 OR keep_alive = 1 OR start_on_mount = 1 OR program LIKE '/tmp/%' OR program LIKE '%/Users/%' OR (program NOT LIKE '/usr/%' AND program NOT LIKE '/System/%' AND program NOT LIKE '/Library/Application Support/Apple/%') OR inetd_compatibility = 'true') ORDER BY disabled ASC, inetd_compatibility DESC, run_at_load DESC;",
        "updated_at": "2025-01-30T12:00:00.000Z",
        "updated_by": "elastic"
    },
    "coreMigrationVersion": "8.3.0",
    "id": "osquery_manager-a1b2c3d4-5678-4abc-def0-123456789036",
    "references": [],
    "type": "osquery-saved-query",
    "updated_at": "2025-01-30T12:00:00.000Z",
    "version": "WzEsMV0="
}
