# packages/osquery_manager/kibana/osquery_saved_query/osquery_manager-f8b0894b-772d-4242-8e19-dbc5d7ae2e06.json
{
  "attributes": {
    "created_at": "2025-11-06T00:00:00.000Z",
    "created_by": "elastic",
    "description": "Detects suspicious systemd persistence on Linux with multi-factor risk scoring. Monitors /etc/systemd/system/ (80% of malware), /etc/systemd/user/, /run/systemd/system/, ~/.config/systemd/user/, /tmp, and /var/tmp. Prioritizes socket/timer units, enabled services, and recently created units. Risk levels: CRITICAL (70+), HIGH (50-69), MEDIUM (30-49), LOW (<30). Returns entries with risk_score >= 30 (MEDIUM+) or created within last 7 days.",
    "ecs_mapping": [
      {
        "key": "service.name",
        "value": {
          "field": "id"
        }
      },
      {
        "key": "service.state",
        "value": {
          "field": "active_state"
        }
      },
      {
        "key": "file.path",
        "value": {
          "field": "fragment_path"
        }
      },
      {
        "key": "event.category",
        "value": {
          "value": [
            "configuration"
          ]
        }
      },
      {
        "key": "event.type",
        "value": {
          "value": [
            "info"
          ]
        }
      },
      {
        "key": "tags",
        "value": {
          "value": [
            "persistence",
            "systemd",
            "linux"
          ]
        }
      }
    ],
    "id": "suspicious_systemd_linux_elastic",
    "interval": "3600",
    "platform": "linux",
    "query": "WITH systemd_base AS (\n  SELECT\n    su.id,\n    su.description,\n    su.load_state,\n    su.active_state,\n    su.sub_state,\n    su.fragment_path,\n    su.unit_file_state,\n    f.ctime,\n    f.mtime,\n    CASE\n      WHEN su.id LIKE '%.socket' THEN 'SOCKET'\n      WHEN su.id LIKE '%.timer' THEN 'TIMER'\n      WHEN su.id LIKE '%.service' THEN 'SERVICE'\n      WHEN su.id LIKE '%.target' THEN 'TARGET'\n      WHEN su.id LIKE '%.path' THEN 'PATH'\n      ELSE 'OTHER'\n    END AS unit_type\n  FROM systemd_units su\n  LEFT JOIN file f ON f.path = su.fragment_path\n  WHERE (\n    su.fragment_path LIKE '/etc/systemd/system/%'\n    OR su.fragment_path LIKE '/etc/systemd/user/%'\n    OR su.fragment_path LIKE '/run/systemd/system/%'\n    OR su.fragment_path LIKE '%/.config/systemd/user/%'\n    OR su.fragment_path LIKE '/tmp/%'\n    OR su.fragment_path LIKE '/var/tmp/%'\n  )\n  AND su.fragment_path NOT LIKE '/lib/systemd/%'\n  AND su.fragment_path NOT LIKE '/usr/lib/systemd/%'\n),\nsystemd_risk AS (\n  SELECT\n    id,\n    description,\n    load_state,\n    active_state,\n    sub_state,\n    fragment_path,\n    unit_file_state,\n    unit_type,\n    ctime,\n    mtime,\n    CASE\n      WHEN fragment_path LIKE '/tmp/%' OR fragment_path LIKE '/var/tmp/%' THEN 40\n      WHEN fragment_path LIKE '/etc/systemd/system/%' THEN 30\n      WHEN fragment_path LIKE '%/.config/systemd/user/%' THEN 25\n      WHEN fragment_path LIKE '/run/systemd/system/%' THEN 20\n      WHEN fragment_path LIKE '/etc/systemd/user/%' THEN 15\n      ELSE 0\n    END AS location_risk,\n    CASE\n      WHEN unit_file_state = 'enabled' AND active_state = 'active' THEN 35\n      WHEN unit_file_state = 'enabled' AND active_state = 'inactive' THEN 25\n      WHEN unit_file_state = 'static' AND active_state = 'active' THEN 20\n      WHEN load_state = 'error' THEN 15\n      WHEN load_state = 'masked' THEN 0\n      ELSE 10\n    END AS state_risk,\n    CASE\n      WHEN unit_type = 'SOCKET' THEN 25\n      WHEN unit_type = 'TIMER' THEN 20\n      WHEN unit_type = 'SERVICE' THEN 10\n      ELSE 0\n    END AS type_risk\n  FROM systemd_base\n)\nSELECT\n  id,\n  description,\n  fragment_path,\n  unit_file_state,\n  active_state,\n  load_state,\n  unit_type,\n  location_risk,\n  state_risk,\n  type_risk,\n  (location_risk + state_risk + type_risk) AS risk_score,\n  CASE\n    WHEN (location_risk + state_risk + type_risk) >= 70 THEN 'CRITICAL'\n    WHEN (location_risk + state_risk + type_risk) >= 50 THEN 'HIGH'\n    WHEN (location_risk + state_risk + type_risk) >= 30 THEN 'MEDIUM'\n    ELSE 'LOW'\n  END AS risk_level,\n  CASE\n    WHEN unit_type = 'SOCKET' THEN 'HIGH: Socket unit detected - potential network backdoor listener'\n    WHEN unit_type = 'TIMER' THEN 'MEDIUM: Timer unit detected - potential C2 beaconing mechanism'\n    WHEN fragment_path LIKE '/tmp/%' OR fragment_path LIKE '/var/tmp/%' THEN 'CRITICAL: Unit file in temporary directory - likely malware persistence'\n    WHEN fragment_path LIKE '/etc/systemd/system/%' AND unit_file_state = 'enabled' THEN 'HIGH: Enabled system unit - review configuration for suspicious commands'\n    WHEN unit_file_state = 'enabled' AND active_state = 'active' THEN 'MEDIUM: Enabled and active unit - validate legitimacy'\n    ELSE 'Review unit configuration and validate against baseline'\n  END AS analyst_guidance,\n  datetime(ctime, 'unixepoch', 'localtime') AS file_created,\n  datetime(mtime, 'unixepoch', 'localtime') AS file_modified,\n  CAST((strftime('%s', 'now') - COALESCE(mtime, strftime('%s', 'now'))) / 86400 AS INTEGER) AS age_days\nFROM systemd_risk\nWHERE (\n  (location_risk + state_risk + type_risk) >= 30\n  OR CAST((strftime('%s', 'now') - COALESCE(mtime, strftime('%s', 'now'))) / 86400 AS INTEGER) <= 7\n  OR unit_type IN ('SOCKET', 'TIMER')\n)\nORDER BY (location_risk + state_risk + type_risk) DESC, CAST((strftime('%s', 'now') - COALESCE(mtime, strftime('%s', 'now'))) / 86400 AS INTEGER) ASC, id\nLIMIT 200;",
    "updated_at": "2025-11-06T00:00:00.000Z",
    "updated_by": "elastic"
  },
  "coreMigrationVersion": "9.2.0",
  "id": "osquery_manager-f8b0894b-772d-4242-8e19-dbc5d7ae2e06",
  "references": [],
  "type": "osquery-saved-query",
  "updated_at": "2025-11-05T09:19:01.000Z",
  "version": "Wzk2LDJd"
}
