# packages/osquery_manager/kibana/osquery_saved_query/osquery_manager-5823a22e-5add-416d-a142-de323400edb0.json
{
  "attributes": {
    "created_at": "2025-01-06T12:00:00.000Z",
    "created_by": "elastic",
    "description": "Detects suspicious third-party macOS launchd persistence. Excludes Apple services (/System/, /Library/Apple/, Apple-signed) using regex. Focuses on: unsigned binaries, programs in /tmp or user directories, hidden files, third-party services with run_at_load. Includes signature validation, hash enrichment, and service_type classification. Optimized for low noise (~10-20 results vs 600+). MITRE ATT&CK: T1543.001 (Launch Agent), T1543.004 (Launch Daemon).",
    "ecs_mapping": [
      {
        "key": "service.name",
        "value": {
          "field": "label"
        }
      },
      {
        "key": "file.path",
        "value": {
          "field": "path"
        }
      },
      {
        "key": "process.executable",
        "value": {
          "field": "program"
        }
      },
      {
        "key": "user.name",
        "value": {
          "field": "username"
        }
      },
      {
        "key": "file.hash.sha256",
        "value": {
          "field": "sha256"
        }
      },
      {
        "key": "file.code_signature.subject_name",
        "value": {
          "field": "signer"
        }
      },
      {
        "key": "file.code_signature.status",
        "value": {
          "field": "signature_status"
        }
      },
      {
        "key": "service.type",
        "value": {
          "field": "service_type"
        }
      },
      {
        "key": "event.category",
        "value": {
          "value": [
            "configuration",
            "process"
          ]
        }
      },
      {
        "key": "event.type",
        "value": {
          "value": [
            "info"
          ]
        }
      },
      {
        "key": "event.kind",
        "value": {
          "value": [
            "state"
          ]
        }
      },
      {
        "key": "tags",
        "value": {
          "value": [
            "threat_hunting",
            "suspicious",
            "darwin",
            "macos",
            "launchd",
            "services",
            "persistence"
          ]
        }
      }
    ],
    "id": "services_suspicious_darwin",
    "interval": "3600",
    "timeout": 120,
    "platform": "darwin",
    "query": "-- macOS Suspicious LaunchD Services Detection (Optimized)\n-- Source: Native osquery tables (launchd, signature, hash)\n-- Focus: Third-party persistence only - excludes all Apple services\n-- MITRE ATT&CK: T1543.001 (Launch Agent), T1543.004 (Launch Daemon)\n\nWITH services_enriched AS (\n  SELECT\n    l.label,\n    l.name,\n    l.path,\n    l.program,\n    l.program_arguments,\n    l.run_at_load,\n    l.keep_alive,\n    COALESCE(l.username, 'root') AS username,\n    COALESCE(s.signed, 0) AS is_signed,\n    s.authority AS signer,\n    CASE WHEN s.signed = 1 THEN 'signed' ELSE 'unsigned' END AS signature_status,\n    (SELECT sha256 FROM hash WHERE path = l.program LIMIT 1) AS sha256,\n    -- Classify service type for analysis\n    CASE\n      WHEN l.path LIKE '/Library/LaunchDaemons/%' THEN 'System-Daemon'\n      WHEN l.path LIKE '/Library/LaunchAgents/%' THEN 'System-Agent'\n      WHEN l.path LIKE '%/Library/LaunchAgents/%' THEN 'User-Agent'\n      ELSE 'Other'\n    END AS service_type\n  FROM launchd l\n  LEFT JOIN signature s ON s.path = l.program\n  WHERE l.path IS NOT NULL AND l.path != ''\n)\n\nSELECT *\nFROM services_enriched\n-- EXCLUDE all Apple/system services using consolidated regex\nWHERE regex_match(\n    path,\n    '^(/System/|/Library/Apple/|/usr/libexec/)',\n    0\n) IS NULL\n  AND regex_match(\n    COALESCE(signer, ''),\n    '(Apple|Software Signing)',\n    0\n) IS NULL\n-- Focus on suspicious indicators\nAND (\n  -- Unsigned is always suspicious\n  is_signed = 0\n  -- Or in suspicious locations\n  OR program LIKE '/Users/%'\n  OR program LIKE '/tmp/%'\n  OR program LIKE '/var/tmp/%'\n  OR program LIKE '/private/tmp/%'\n  OR program LIKE '%/.%'\n  -- Or third-party with persistence enabled\n  OR (is_signed = 1 AND run_at_load = '1')\n)\nORDER BY\n  is_signed ASC,\n  service_type,\n  label;",
    "updated_at": "2025-11-28T13:50:00.000Z",
    "updated_by": "elastic"
  },
  "coreMigrationVersion": "9.2.0",
  "id": "osquery_manager-5823a22e-5add-416d-a142-de323400edb0",
  "references": [],
  "type": "osquery-saved-query",
  "updated_at": "2025-11-28T13:50:00.000Z",
  "version": "Wzk1LDJd"
}
