# packages/osquery_manager/kibana/osquery_saved_query/osquery_manager-ec67184a-3a16-4b6e-bb75-5575b7b1ec83.json
{
  "attributes": {
    "created_at": "2025-11-05T09:19:01.000Z",
    "created_by": "elastic",
    "description": "Detect suspicious startup items on macOS, excluding system paths but flagging user directories and temporary locations",
    "ecs_mapping": [
      {
        "key": "process.name",
        "value": {
          "field": "name"
        }
      },
      {
        "key": "process.executable",
        "value": {
          "field": "path"
        }
      },
      {
        "key": "process.args",
        "value": {
          "field": "args"
        }
      },
      {
        "key": "registry.path",
        "value": {
          "field": "source"
        }
      },
      {
        "key": "user.name",
        "value": {
          "field": "username"
        }
      },
      {
        "key": "event.category",
        "value": {
          "value": [
            "process"
          ]
        }
      },
      {
        "key": "event.type",
        "value": {
          "value": [
            "info"
          ]
        }
      },
      {
        "key": "tags",
        "value": {
          "value": [
            "persistence",
            "startup-items",
            "macos"
          ]
        }
      }
    ],
    "id": "startup_items_persistence_macos_elastic",
    "interval": "3600",
    "platform": "darwin",
    "query": "SELECT \n  name, \n  path, \n  args, \n  type, \n  source, \n  status, \n  username \nFROM startup_items \nWHERE path NOT LIKE '/System/%' \n  AND path NOT LIKE '/Applications/Xcode.app/%' \n  AND path NOT LIKE '/Library/Apple/%' \n  AND name NOT LIKE 'com.apple.%' \n  AND (\n    path LIKE '/tmp/%' \n    OR path LIKE '/var/tmp/%' \n    OR path LIKE '/Users/%/Library/Application Support/%' \n    OR path LIKE '%/Downloads/%' \n    OR path LIKE '%/Desktop/%'\n  ) \nORDER BY name \nLIMIT 100;",
    "updated_at": "2025-11-05T09:19:01.000Z",
    "updated_by": "elastic"
  },
  "coreMigrationVersion": "9.3.0",
  "id": "osquery_manager-ec67184a-3a16-4b6e-bb75-5575b7b1ec83",
  "references": [],
  "type": "osquery-saved-query",
  "updated_at": "2025-11-05T09:19:01.000Z",
  "version": "Wzk4LDJd"
}
