# packages/osquery_manager/kibana/osquery_saved_query/osquery_manager-265051dd-bc20-491a-a998-98ebc2f00af7.json
{
    "attributes": {
        "created_at": "2025-12-08T00:00:00.000Z",
        "created_by": "elastic",
        "description": "Detects ALL suspicious Windows scheduled tasks using dual-detection approach for complete coverage: (1) Non-whitelisted tasks not in known-good allowlist and (2) Living off the Land (LOTL) attack indicators. This query intentionally excludes hash enrichment to capture tasks even when executable files are orphaned, deleted, or inaccessible. Identifies abuse of legitimate Windows tools (powershell -e, certutil, wscript, mshta, etc.) for persistence, privilege escalation, and lateral movement. Uses comprehensive path-based whitelist filtering to auto-exclude legitimate Microsoft Windows tasks from System32/SysWOW64, Windows subsystems (Bluetooth, UPnP, Diagnosis, etc.), and known-good third-party updaters (Google, Adobe, OneDrive). For tasks with existing files requiring hash/signature analysis, use the companion query 'scheduled_tasks_enriched_windows_elastic'. Maps to MITRE ATT&CK T1053.005 (Scheduled Task), T1059.001 (PowerShell), T1105 (Ingress Tool Transfer).",
        "ecs_mapping": [
            {
                "key": "event.category",
                "value": {
                    "value": ["configuration"]
                }
            },
            {
                "key": "event.type",
                "value": {
                    "value": ["info"]
                }
            },
            {
                "key": "event.action",
                "value": {
                    "value": "osquery.scheduled_tasks"
                }
            },
            {
                "key": "service.name",
                "value": {
                    "field": "task_name"
                }
            },
            {
                "key": "service.state",
                "value": {
                    "field": "task_state"
                }
            },
            {
                "key": "process.command_line",
                "value": {
                    "field": "command_line"
                }
            },
            {
                "key": "file.path",
                "value": {
                    "field": "task_path"
                }
            },
            {
                "key": "event.created",
                "value": {
                    "field": "last_execution"
                }
            },
            {
                "key": "rule.category",
                "value": {
                    "field": "task_type"
                }
            },
            {
                "key": "rule.description",
                "value": {
                    "field": "detection_reason"
                }
            },
            {
                "key": "rule.name",
                "value": {
                    "field": "detection_method"
                }
            },
            {
                "key": "tags",
                "value": {
                    "value": ["osquery", "persistence", "scheduled_tasks", "windows"]
                }
            }
        ],
        "id": "scheduled_tasks_suspicious_windows_elastic",
        "interval": "3600",
        "platform": "windows",
        "query": "-- Detects ALL suspicious Windows scheduled tasks (full coverage, no hash enrichment)\n-- MITRE: T1053.005, T1059.001, T1105\n-- Detection: NON_WHITELISTED + LOTL_INDICATOR dual approach\n-- Use companion query 'scheduled_tasks_enriched_windows_elastic' for hash/signature analysis\n\nWITH non_whitelisted AS (\n    SELECT \n        st.name AS task_name,\n        st.action AS command_line,\n        st.path AS task_path,\n        st.enabled,\n        st.state AS task_state,\n        st.hidden,\n        datetime(st.last_run_time, 'unixepoch') AS last_execution,\n        datetime(st.next_run_time, 'unixepoch') AS next_execution,\n        st.last_run_message,\n        1 AS is_non_whitelisted,\n        0 AS is_lotl\n    FROM scheduled_tasks AS st\n    WHERE st.name IS NOT NULL\n        AND st.name != ''\n        AND st.action IS NOT NULL\n        AND st.action != ''\n        -- Allowlist: Microsoft system tasks from System32/SysWOW64\n        AND NOT (\n            (st.path LIKE '\\Microsoft\\Windows\\%' OR st.path LIKE '\\Microsoft\\Office\\%' OR st.path LIKE '\\Microsoft\\XblGameSave\\%')\n            AND (\n                LOWER(st.action) LIKE '%systemroot%\\system32\\%'\n                OR LOWER(st.action) LIKE '%windir%\\system32\\%'\n                OR LOWER(st.action) LIKE 'c:\\windows\\system32\\%'\n                OR LOWER(st.action) LIKE '%systemroot%\\syswow64\\%'\n                OR LOWER(st.action) LIKE '%windir%\\syswow64\\%'\n                OR LOWER(st.action) LIKE 'c:\\windows\\syswow64\\%'\n                OR LOWER(st.action) LIKE '%programfiles%\\windows defender\\%'\n                OR LOWER(st.action) LIKE 'c:\\program files\\windows defender\\%'\n                OR LOWER(st.action) LIKE '%program files%\\microsoft office\\%'\n                OR LOWER(st.action) LIKE 'c:\\program files\\microsoft office\\%'\n                OR LOWER(st.action) LIKE 'c:\\program files (x86)\\microsoft office\\%'\n            )\n        )\n        -- Allowlist: Defender tasks from legitimate ProgramData\n        AND NOT (\n            st.path LIKE '\\Microsoft\\Windows\\Windows Defender\\%'\n            AND (\n                LOWER(st.action) LIKE 'c:\\programdata\\microsoft\\windows defender\\%'\n                OR LOWER(st.action) LIKE '%\\programdata\\microsoft\\windows defender\\%'\n                OR LOWER(st.action) LIKE '%\\program files\\windows defender\\%'\n            )\n        )\n        -- Allowlist: Edge Update tasks\n        AND NOT (\n            st.path LIKE '\\MicrosoftEdgeUpdate%'\n            OR (st.path LIKE '%EdgeUpdate%' AND st.action LIKE '%Microsoft\\EdgeUpdate\\%')\n        )\n        -- Allowlist: Common Windows executables for Microsoft tasks\n        AND NOT (\n            st.path LIKE '\\Microsoft\\Windows\\%'\n            AND (\n                LOWER(st.action) LIKE 'sc.exe %'\n                OR LOWER(st.action) LIKE 'btudtask.exe%'\n                OR LOWER(st.action) LIKE 'schtasks.exe%'\n                OR LOWER(st.action) LIKE '%wmpnscfg.exe%'\n            )\n        )\n        -- Allowlist: Windows subsystem tasks\n        AND NOT (\n            st.path LIKE '%\\Bluetooth\\%'\n            OR st.path LIKE '%\\UPnP\\%'\n            OR st.path LIKE '%\\UsageAndQualityInsights\\%'\n            OR st.path LIKE '%\\Diagnosis\\%'\n            OR st.path LIKE '%\\UpdateOrchestrator\\%'\n            OR st.path LIKE '%\\WindowsUpdate\\%'\n            OR st.path LIKE '%\\Defrag\\%'\n            OR st.path LIKE '%\\Application Experience\\%'\n            OR st.path LIKE '%\\Location\\%'\n            OR st.path LIKE '%\\Sysmain\\%'\n            OR st.path LIKE '%\\ApplicationData\\%'\n            OR st.path LIKE '%\\DUSM\\%'\n            OR st.path LIKE '%\\capabilityaccessmanager\\%'\n            OR st.path LIKE '%\\AppxDeploymentClient\\%'\n            OR st.path LIKE '%\\Windows Media Sharing\\%'\n            OR st.path LIKE '%\\Speech\\%'\n            OR st.path LIKE '%\\TextServicesFramework\\%'\n            OR st.path LIKE '%\\Power Efficiency Diagnostics\\%'\n            OR st.path LIKE '%\\RecoveryEnvironment\\%'\n            OR st.path LIKE '%\\Registry\\%'\n            OR st.path LIKE '%\\MemoryDiagnostic\\%'\n            OR st.path LIKE '%\\Maintenance\\%'\n            OR st.path LIKE '%\\DiskCleanup\\%'\n            OR st.path LIKE '%\\Servicing\\%'\n            OR st.path LIKE '%\\Software Inventory Logging\\%'\n        )\n        -- Allowlist: Known-good third-party updaters (Google, Adobe, OneDrive)\n        AND NOT (\n            st.path LIKE '%\\Google\\Update%'\n            AND st.action LIKE '%Google\\Update\\%'\n        )\n        AND NOT (\n            st.path LIKE '%\\Adobe\\%'\n            AND st.action LIKE '%Adobe%'\n        )\n        AND NOT (\n            st.path LIKE '%OneDrive%'\n        )\n        -- Allowlist: Intel/hardware vendor tasks\n        AND NOT (\n            st.path LIKE '%\\Intel\\%'\n            AND LOWER(st.action) LIKE '%\\intel\\%'\n        )\n),\nlotl_indicators AS (\n    SELECT \n        st.name AS task_name,\n        st.action AS command_line,\n        st.path AS task_path,\n        st.enabled,\n        st.state AS task_state,\n        st.hidden,\n        datetime(st.last_run_time, 'unixepoch') AS last_execution,\n        datetime(st.next_run_time, 'unixepoch') AS next_execution,\n        st.last_run_message,\n        0 AS is_non_whitelisted,\n        1 AS is_lotl,\n        CASE\n            -- PowerShell LOTL patterns\n            WHEN LOWER(st.action) LIKE '%powershell% -e%' \n                OR LOWER(st.action) LIKE '% -enc %' \n                OR LOWER(st.action) LIKE '% -encodedcommand %' \n            THEN 'PowerShell base64 encoded command'\n            WHEN LOWER(st.action) LIKE '%invoke-webrequest%' \n                OR LOWER(st.action) LIKE '%iwr %' \n                OR LOWER(st.action) LIKE '%invoke-restmethod%' \n            THEN 'PowerShell download cradle'\n            WHEN LOWER(st.action) LIKE '%-executionpolicy bypass%' \n                OR LOWER(st.action) LIKE '%-ep bypass%' \n            THEN 'PowerShell execution policy bypass'\n            WHEN LOWER(st.action) LIKE '%-w hidden%' \n                OR LOWER(st.action) LIKE '%-windowstyle hidden%' \n            THEN 'PowerShell hidden window'\n            WHEN LOWER(st.action) LIKE '%iex(%' \n                OR LOWER(st.action) LIKE '%invoke-expression%' \n            THEN 'PowerShell Invoke-Expression abuse'\n            WHEN LOWER(st.action) LIKE '%[convert]::frombase64string%' \n            THEN 'PowerShell base64 decode obfuscation'\n            WHEN LOWER(st.action) LIKE '%downloadstring%' \n                OR LOWER(st.action) LIKE '%downloadfile%' \n            THEN 'PowerShell remote download method'\n            WHEN LOWER(st.action) LIKE '%new-object net.webclient%' \n                OR LOWER(st.action) LIKE '%system.net.webclient%' \n            THEN 'PowerShell WebClient abuse'\n            -- Download/transfer utilities\n            WHEN LOWER(st.action) LIKE '%certutil% -urlcache%' \n                OR LOWER(st.action) LIKE '%certutil% -f%' \n            THEN 'CertUtil download/decode abuse'\n            WHEN LOWER(st.action) LIKE '%bitsadmin% /transfer%' \n                OR LOWER(st.action) LIKE '%bitsadmin%/transfer%' \n            THEN 'BITSAdmin download abuse'\n            -- Suspicious execution paths\n            WHEN (LOWER(st.action) LIKE '%c:\\users\\public\\%' \n                OR LOWER(st.action) LIKE '%c:\\programdata\\%') \n                AND LOWER(st.action) NOT LIKE '%\\programdata\\microsoft\\windows defender\\%'\n                AND LOWER(st.action) NOT LIKE '%\\programdata\\microsoft\\%' \n            THEN 'Suspicious file path (writable by low-priv users)'\n            WHEN LOWER(st.action) LIKE '%\\temp\\%' \n                OR LOWER(st.action) LIKE '%\\appdata\\local\\temp\\%' \n            THEN 'Execution from Temp directory'\n            -- Script host abuse\n            WHEN LOWER(st.action) LIKE '%wscript.exe%' \n                OR LOWER(st.action) LIKE '%cscript.exe%' \n            THEN 'Windows Script Host abuse'\n            WHEN LOWER(st.action) LIKE '%mshta.exe%' \n                OR LOWER(st.action) LIKE '%mshta %' \n            THEN 'MSHTA.exe abuse'\n            -- Script file execution\n            WHEN LOWER(st.action) LIKE '%.hta%' \n                OR LOWER(st.action) LIKE '%.vbs%' \n                OR LOWER(st.action) LIKE '%.vbe%' \n                OR LOWER(st.action) LIKE '%.js%' \n                OR LOWER(st.action) LIKE '%.jse%' \n                OR LOWER(st.action) LIKE '%.wsf%' \n            THEN 'Script file execution (.hta/.vbs/.js)'\n            -- Proxy execution\n            WHEN (LOWER(st.action) LIKE '%regsvr32%' OR LOWER(st.action) LIKE '%rundll32%') \n                AND st.path NOT LIKE '\\Microsoft\\Windows\\%' \n            THEN 'Proxy execution via regsvr32/rundll32'\n            WHEN LOWER(st.action) LIKE '%msiexec%' \n                AND st.path NOT LIKE '\\Microsoft\\Windows\\%' \n            THEN 'MSI package execution'\n            -- Hidden task detection\n            WHEN st.hidden = 1 AND st.enabled = 1\n                AND st.path NOT LIKE '\\Microsoft\\Windows\\%'\n                AND st.path NOT LIKE '\\Microsoft\\Office\\%'\n                AND st.path NOT LIKE '\\Microsoft\\XblGameSave\\%'\n                AND st.path NOT LIKE '\\Microsoft\\EdgeUpdate\\%'\n                AND st.path NOT LIKE '\\MicrosoftEdgeUpdate%'\n                AND st.path NOT LIKE '%OneDrive%'\n            THEN 'Hidden and enabled scheduled task'\n            -- Network utilities\n            WHEN LOWER(st.action) LIKE '%net.exe %' \n                AND (LOWER(st.action) LIKE '% user %' OR LOWER(st.action) LIKE '% localgroup %')\n            THEN 'Net.exe user/group enumeration'\n            WHEN LOWER(st.action) LIKE '%nltest%' \n            THEN 'Domain trust enumeration (nltest)'\n            ELSE 'Unknown LOTL pattern'\n        END AS detection_reason\n    FROM scheduled_tasks AS st\n    WHERE st.name IS NOT NULL\n        AND st.name != ''\n        AND st.action IS NOT NULL\n        AND st.action != ''\n        AND (\n            -- PowerShell patterns\n            LOWER(st.action) LIKE '%powershell% -e%'\n            OR LOWER(st.action) LIKE '% -enc %'\n            OR LOWER(st.action) LIKE '% -encodedcommand %'\n            OR LOWER(st.action) LIKE '%invoke-webrequest%'\n            OR LOWER(st.action) LIKE '%iwr %'\n            OR LOWER(st.action) LIKE '%invoke-restmethod%'\n            OR LOWER(st.action) LIKE '%-executionpolicy bypass%'\n            OR LOWER(st.action) LIKE '%-ep bypass%'\n            OR LOWER(st.action) LIKE '%-w hidden%'\n            OR LOWER(st.action) LIKE '%-windowstyle hidden%'\n            OR LOWER(st.action) LIKE '%iex(%'\n            OR LOWER(st.action) LIKE '%invoke-expression%'\n            OR LOWER(st.action) LIKE '%[convert]::frombase64string%'\n            OR LOWER(st.action) LIKE '%downloadstring%'\n            OR LOWER(st.action) LIKE '%downloadfile%'\n            OR LOWER(st.action) LIKE '%new-object net.webclient%'\n            OR LOWER(st.action) LIKE '%system.net.webclient%'\n            -- Download utilities\n            OR LOWER(st.action) LIKE '%certutil% -urlcache%'\n            OR LOWER(st.action) LIKE '%certutil% -f%'\n            OR LOWER(st.action) LIKE '%bitsadmin% /transfer%'\n            OR LOWER(st.action) LIKE '%bitsadmin%/transfer%'\n            -- Suspicious paths\n            OR (LOWER(st.action) LIKE '%c:\\users\\public\\%' AND LOWER(st.action) NOT LIKE '%\\programdata\\microsoft\\windows defender\\%')\n            OR (LOWER(st.action) LIKE '%c:\\programdata\\%' AND LOWER(st.action) NOT LIKE '%\\programdata\\microsoft\\windows defender\\%' AND LOWER(st.action) NOT LIKE '%\\programdata\\microsoft\\%')\n            OR LOWER(st.action) LIKE '%\\temp\\%'\n            OR LOWER(st.action) LIKE '%\\appdata\\local\\temp\\%'\n            -- Script hosts\n            OR LOWER(st.action) LIKE '%wscript.exe%'\n            OR LOWER(st.action) LIKE '%cscript.exe%'\n            OR LOWER(st.action) LIKE '%mshta.exe%'\n            OR LOWER(st.action) LIKE '%mshta %'\n            -- Script files\n            OR LOWER(st.action) LIKE '%.hta%'\n            OR LOWER(st.action) LIKE '%.vbs%'\n            OR LOWER(st.action) LIKE '%.vbe%'\n            OR LOWER(st.action) LIKE '%.js%'\n            OR LOWER(st.action) LIKE '%.jse%'\n            OR LOWER(st.action) LIKE '%.wsf%'\n            -- Proxy execution\n            OR (LOWER(st.action) LIKE '%regsvr32%' AND st.path NOT LIKE '\\Microsoft\\Windows\\%')\n            OR (LOWER(st.action) LIKE '%rundll32%' AND st.path NOT LIKE '\\Microsoft\\Windows\\%')\n            OR (LOWER(st.action) LIKE '%msiexec%' AND st.path NOT LIKE '\\Microsoft\\Windows\\%')\n            -- Hidden tasks\n            OR (\n                st.hidden = 1 AND st.enabled = 1\n                AND st.path NOT LIKE '\\Microsoft\\Windows\\%'\n                AND st.path NOT LIKE '\\Microsoft\\Office\\%'\n                AND st.path NOT LIKE '\\Microsoft\\XblGameSave\\%'\n                AND st.path NOT LIKE '\\Microsoft\\EdgeUpdate\\%'\n                AND st.path NOT LIKE '\\MicrosoftEdgeUpdate%'\n                AND st.path NOT LIKE '%OneDrive%'\n            )\n            -- Network utilities\n            OR (LOWER(st.action) LIKE '%net.exe %' AND (LOWER(st.action) LIKE '% user %' OR LOWER(st.action) LIKE '% localgroup %'))\n            OR LOWER(st.action) LIKE '%nltest%'\n        )\n),\ncombined AS (\n    SELECT \n        st.name AS task_name,\n        st.action AS command_line,\n        st.path AS task_path,\n        st.enabled,\n        st.state AS task_state,\n        st.hidden,\n        datetime(st.last_run_time, 'unixepoch') AS last_execution,\n        datetime(st.next_run_time, 'unixepoch') AS next_execution,\n        st.last_run_message,\n        CASE \n            WHEN MAX(COALESCE(li.is_lotl, 0)) = 1 \n            THEN 'Scheduled Task (LOTL)'\n            ELSE 'Scheduled Task'\n        END AS task_type,\n        CASE \n            WHEN MAX(COALESCE(li.is_lotl, 0)) = 1 AND MAX(COALESCE(nw.is_non_whitelisted, 0)) = 1 \n            THEN 'LOTL_INDICATOR + NON_WHITELISTED'\n            WHEN MAX(COALESCE(li.is_lotl, 0)) = 1 \n            THEN 'LOTL_INDICATOR'\n            ELSE 'NON_WHITELISTED'\n        END AS detection_method,\n        CASE \n            WHEN MAX(COALESCE(li.is_lotl, 0)) = 1 AND MAX(COALESCE(nw.is_non_whitelisted, 0)) = 1 \n            THEN MAX(li.detection_reason) || ' + Not in known-good allowlist'\n            WHEN MAX(COALESCE(li.is_lotl, 0)) = 1 \n            THEN MAX(li.detection_reason)\n            ELSE 'Scheduled task not in known-good allowlist'\n        END AS detection_reason\n    FROM scheduled_tasks AS st\n    LEFT JOIN non_whitelisted AS nw ON st.name = nw.task_name AND st.path = nw.task_path AND st.action = nw.command_line\n    LEFT JOIN lotl_indicators AS li ON st.name = li.task_name AND st.path = li.task_path AND st.action = li.command_line\n    WHERE COALESCE(nw.is_non_whitelisted, 0) = 1 OR COALESCE(li.is_lotl, 0) = 1\n    GROUP BY st.name, st.path, st.action, st.enabled, st.state, st.hidden, st.last_run_time, st.next_run_time, st.last_run_message\n)\nSELECT \n    c.task_name,\n    c.command_line,\n    c.task_path,\n    c.task_type,\n    c.enabled,\n    c.task_state,\n    c.hidden,\n    c.last_execution,\n    c.next_execution,\n    c.last_run_message,\n    c.detection_method,\n    c.detection_reason\nFROM combined AS c\nORDER BY \n    CASE WHEN c.detection_method LIKE 'LOTL_INDICATOR%' THEN 0 ELSE 1 END,\n    c.detection_reason,\n    c.hidden DESC,\n    c.task_name",
        "updated_at": "2025-12-08T00:00:00.000Z",
        "updated_by": "elastic"
    },
    "coreMigrationVersion": "9.2.0",
    "id": "osquery_manager-265051dd-bc20-491a-a998-98ebc2f00af7",
    "references": [],
    "type": "osquery-saved-query",
    "updated_at": "2025-12-08T00:00:00.000Z",
    "version": "WzEwMDAwLDFd"
}
