# packages/osquery_manager/kibana/osquery_saved_query/osquery_manager-a1b2c3d4-5678-4abc-def0-123456789034.json
{
    "attributes": {
        "created_at": "2025-01-30T12:00:00.000Z",
        "created_by": "elastic",
        "description": "Detect Linux service-based persistence and privilege escalation by monitoring systemd units from suspicious locations, user-level services, and services with unusual configurations",
        "ecs_mapping": [
            {
                "key": "service.name",
                "value": {
                    "field": "id"
                }
            },
            {
                "key": "service.state",
                "value": {
                    "field": "active_state"
                }
            },
            {
                "key": "file.path",
                "value": {
                    "field": "fragment_path"
                }
            },
            {
                "key": "event.category",
                "value": {
                    "value": ["process"]
                }
            },
            {
                "key": "event.type",
                "value": {
                    "value": ["info"]
                }
            },
            {
                "key": "event.action",
                "value": {
                    "value": "systemd-service-enumeration"
                }
            },
            {
                "key": "event.dataset",
                "value": {
                    "value": "linux.systemd"
                }
            },
            {
                "key": "tags",
                "value": {
                    "value": ["persistence", "linux", "systemd", "service", "T1543.002"]
                }
            }
        ],
        "id": "linux_systemd_services_persistence_elastic",
        "interval": "3600",
        "platform": "linux",
        "query": "SELECT id, description, load_state, active_state, sub_state, unit_file_state, fragment_path, source_path FROM systemd_units WHERE type = 'service' AND (active_state = 'active' OR unit_file_state = 'enabled') AND (fragment_path LIKE '/tmp/%' OR fragment_path LIKE '%/home/%' OR fragment_path LIKE '/dev/shm/%' OR fragment_path LIKE '%/.config/systemd/%' OR description = '' OR description LIKE '%test%') ORDER BY active_state DESC, unit_file_state;",
        "updated_at": "2025-01-30T12:00:00.000Z",
        "updated_by": "elastic"
    },
    "coreMigrationVersion": "8.3.0",
    "id": "osquery_manager-a1b2c3d4-5678-4abc-def0-123456789034",
    "references": [],
    "type": "osquery-saved-query",
    "updated_at": "2025-01-30T12:00:00.000Z",
    "version": "WzEsMV0="
}
