# packages/carbon_black_cloud/_dev/deploy/docker/files/config.yml
rules:
  - path: /appservices/v6/orgs/xxxx/alerts/_search
    methods: ["POST"]
    responses:
      - status_code: 200
        body: |
          {"results":[{"type":"DEVICE_CONTROL","id":"test1","legacy_alert_id":"C8EB7306-AF26-4A9A-B677-814B3AF69720","org_key":"ABCD6X3T","create_time":"2020-11-17T22:05:13Z","last_update_time":"2020-11-17T22:05:13Z","first_event_time":"2020-11-17T22:02:16Z","last_event_time":"2020-11-17T22:02:16Z","threat_id":"t5678","severity":3,"category":"WARNING","device_id":2,"device_os":"WINDOWS","device_os_version":"Windows 10 x64","device_name":"DESKTOP-002","device_username":"test34@demo.com","policy_id":6997287,"policy_name":"Standard","target_value":"MEDIUM","workflow":{"state":"OPEN","remediation":"","last_update_time":"2020-11-17T22:02:16Z","comment":"","changed_by":"Carbon Black"},"reason":"Access attempted on unapproved USB device SanDisk U3 Cruzer Micro (SN: 0875920EF7C2A304). A Deny Policy Action was applied.","reason_code":"6D578342-9DE5-4353-9C25-1D3D857BFC5B:DCAEB1FA-513C-4026-9AB6-37A935873FBC","device_location":"UNKNOWN","threat_cause_threat_category":"NON_MALWARE","threat_cause_vector":"REMOVABLE_MEDIA","threat_cause_cause_event_id":"FCEE2AF0-D832-4C9F-B988-F11B46028C9E","sensor_action":"DENY","run_state":"DID_NOT_RUN","policy_applied":"APPLIED","vendor_name":"SanDisk","vendor_id":"0x0781","product_name":"U3 Cruzer Micro","product_id":"0x5406","serial_number":"0875920EF7C2A304"}],"num_found":6197,"num_available":6197}
  - path: /api/alerts/v7/orgs/xxxx/alerts/_search
    methods: ["POST"]
    responses:
      - status_code: 200
        body: |
          {"results":[{"org_key":"7DESJ9GN","alert_url":"defense.conferdeploy.net/alerts?s[c][query_string]=id:1c6aba68-24cc-41e3-ad8e-4b545a587b55&orgKey=7DESJ9GN","id":"1c6aba68-24cc-41e3-ad8e-4b545a587b55","type":"WATCHLIST","backend_timestamp":"2024-03-13T08:03:29.540Z","user_update_timestamp":null,"backend_update_timestamp":"2024-03-13T08:03:29.540Z","detection_timestamp":"2024-03-13T08:02:36.578Z","first_event_timestamp":"2024-03-13T08:00:09.894Z","last_event_timestamp":"2024-03-13T08:00:09.894Z","severity":5,"reason":"Process powershell.exe was detected by the report \"Execution - AMSI - New Fileless Scheduled Task Behavior Detected\" in watchlist \"AMSI Threat Intelligence\"","reason_code":"c21ca826-573a-3d97-8c1e-93c8471aab7f:8033b29d-81d2-3c47-82d2-f4a7f398b85d","threat_id":"C21CA826573A8D974C1E93C8471AAB7F","primary_event_id":"re9M9hp8TbGLqyk6QXqQqA-0","policy_applied":"NOT_APPLIED","run_state":"RAN","sensor_action":"ALLOW","workflow":{"change_timestamp":"2024-03-13T08:03:29.540Z","changed_by_type":"SYSTEM","changed_by":"ALERT_CREATION","closure_reason":"NO_REASON","status":"OPEN"},"determination":{"change_timestamp":"2024-03-13T08:03:29.540Z","value":"NONE","changed_by_type":"SYSTEM","changed_by":"ALERT_CREATION"},"tags":null,"alert_notes_present":false,"threat_notes_present":false,"is_updated":false,"device_id":6612391,"device_name":"EIP\\WW-20002","device_uem_id":"","device_target_value":"MEDIUM","device_policy":"default","device_policy_id":6525,"device_os":"WINDOWS","device_os_version":"Windows 10 x64","device_username":"EIP\\Administrator","device_location":"UNKNOWN","device_external_ip":"75.98.230.194","device_internal_ip":"172.16.100.140","asset_group":[],"mdr_alert":false,"mdr_alert_notes_present":false,"mdr_threat_notes_present":false,"report_id":"LrKOC7DtQbm4g8w0UFruQg-d1080521-e617-4e45-94e0-7a145c62c90a","report_name":"Execution - AMSI - New Fileless Scheduled Task Behavior Detected","report_description":"Newer Powershell versions introduced built-in cmdlets to manage scheduled tasks natively without calling out to typical scheduled task processes like at.exe or schtasks.exe. This detection looks for behaviors related to the fileless execution of scheduled tasks. If you are responding to this alert, be sure to correlate the fileless scriptload events with events typically found in your environment Generally, attackers will create scheduled tasks with binaries that are located in user writable directories like AppData, Temp, or public folders.","report_tags":["execution","privesc","persistence","t1053","windows","amsi","attack","attackframework"],"report_link":"https://attack.mitre.org/techniques/T1053/","ioc_id":"d1080521-e617-4e45-94e0-7a145c62c90a","ioc_hit":"(fileless_scriptload_cmdline:Register-ScheduledTask OR fileless_scriptload_cmdline:New-ScheduledTask OR scriptload_content:Register-ScheduledTask OR scriptload_content:New-ScheduledTask) AND NOT (process_cmdline:windows\\\\ccm\\\\systemtemp OR crossproc_name:windows\\\\ccm\\\\ccmexec.exe OR (process_publisher:\"VMware, Inc.\" AND process_publisher_state:FILE_SIGNATURE_STATE_TRUSTED))","watchlists":[{"id":"Ci7w5B4URg6HN60hatQMQ","name":"AMSI Threat Intelligence"}],"process_guid":"7DESJ9GN-0064e5a7-00001434-00000000-1da751c7354ebfe","process_pid":5172,"process_name":"c:\\windows\\system32\\windowspowershell\\v1.0\\powershell.exe","process_sha256":"9785001b0dcf755eddb8af294a373c0b87b2498660f724e76c4d53f9c217c7a3","process_md5":"2e5a8590cf6848968fc23de3fa1e25f1","process_effective_reputation":"TRUSTED_WHITE_LIST","process_reputation":"TRUSTED_WHITE_LIST","process_cmdline":"\"C:\\Windows\\System32\\WindowsPowerShell\\v1.0\\powershell.exe\" -EP Bypass \\\\eip.demo\\sysvol\\EIP.DEMO\\scripts\\Luminol.ps1","process_username":"NT AUTHORITY\\SYSTEM","process_issuer":["Microsoft Windows Production PCA 2011"],"process_publisher":["Microsoft Windows"],"parent_guid":"7DESJ9GN-0064e5a7-0000077c-00000000-1da5ed7ec07b275","parent_pid":1916,"parent_name":"c:\\windows\\system32\\svchost.exe","parent_sha256":"f13de58416730d210dab465b242e9c949fb0a0245eef45b07c381f0c6c8a43c3","parent_md5":"145dcf6706eeea5b066885ee17964c09","parent_effective_reputation":"TRUSTED_WHITE_LIST","parent_reputation":"TRUSTED_WHITE_LIST","parent_cmdline":"C:\\Windows\\system32\\svchost.exe -k netsvcs -p -s Schedule","parent_username":"NT AUTHORITY\\SYSTEM","childproc_guid":"","childproc_username":"","childproc_cmdline":"","ml_classification_final_verdict":"NOT_ANOMALOUS","ml_classification_global_prevalence":"LOW","ml_classification_org_prevalence":"LOW","ml_classification_anomalies":[]}],"num_found":1,"num_available":1}
  - path: /integrationServices/v3/auditlogs
    methods: ["GET"]
    responses:
      - status_code: 200
        body: |
          {"notifications":[{"orgName":"cb-xxxx-xxxx.com","requestUrl":null,"eventTime":1644509070263,"eventId":"2122f8ce8xxxxxxxxxxxxx","loginName":"abc@demo.com","flagged":false,"clientIp":"10.10.10.10","verbose":false,"description":"Logged in successfully"}],"success":true,"message":"Success"}
  - path: /vulnerability/assessment/api/v1/orgs/xxxx/devices/vulnerabilities/summary/_search
    methods: ["POST"]
    responses:
      - status_code: 200
        body: |-
          {"num_found":130,"results":[{"cve_ids":null,"device_id":8,"highest_risk_score":10,"host_name":"DESKTOP-008","last_sync_ts":"2022-01-17T08:33:37.384932Z","name":"DESKTOP-008KK","os_info":{"os_arch":"64-bit","os_name":"Microsoft Windows 10 Education","os_type":"WINDOWS","os_version":"10.0.17763"},"severity":"CRITICAL","sync_status":"COMPLETED","sync_type":"SCHEDULED","type":"ENDPOINT","vm_id":"","vm_name":"","vuln_count":1770}]}
